# systemd/cubie-nas.service
[Unit]
Description=Cubie NAS Management Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/cubie-nas
EnvironmentFile=/opt/cubie-nas/.env
ExecStart=/opt/cubie-nas/.venv/bin/uvicorn app.main:app --host ${APP_HOST} --port ${APP_PORT} --proxy-headers
Restart=always
RestartSec=3
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=false
LimitNOFILE=65535
MemoryMax=300M
CPUAccounting=true
MemoryAccounting=true
TasksMax=256

[Install]
WantedBy=multi-user.target
