# systemd/cubie-nas.service
[Unit]
Description=Cubie NAS Management Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/cubie-nas
EnvironmentFile=/opt/cubie-nas/.env
ExecStart=/opt/cubie-nas/.venv/bin/uvicorn app.main:app --host ${APP_HOST} --port ${APP_PORT} --proxy-headers --workers 1 --no-access-log --timeout-keep-alive 5
Restart=always
RestartSec=3
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/cubie-nas /var/lib/cubie-nas /srv/nas /etc/samba /etc/fstab /etc/cubie-nas /etc/network/interfaces.d
LimitNOFILE=4096
OOMPolicy=continue
MemoryHigh=250M
MemoryMax=300M
CPUAccounting=true
MemoryAccounting=true
TasksMax=256

[Install]
WantedBy=multi-user.target
